# scrumbuddy.cron
00 16 * * 1,2,3,4,5 cd /path/to/scrumbuddy && /usr/bin/rake environment RAILS_ENV=development scrum
09 00 * * 1 cd /path/to/scrumbuddy && /usr/bin/rake environment RAILS_ENV=development weekly
15 00 1 * * cd /path/to/scrumbuddy && /usr/bin/rake environment RAILS_ENV=development monthly
